# assets/config/supervisor/7-entryscript.conf
[program:entryscript]
command=/bin/bash /entry.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/entryscript.log
stderr_logfile=/var/log/entryscript.err